# com.user.tg-cli-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.tg-cli-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/reuben/tg-claude-bridge/tg-bridge</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/reuben/tg-claude-bridge</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/reuben/tg-claude-bridge/logs/tg-bridge.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/reuben/tg-claude-bridge/logs/tg-bridge.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/reuben/.local/bin</string>
    </dict>
</dict>
</plist>
